# deploy/systemd/brc-runtime-monitor.service.d/60-runtime-db-application.conf
[Service]
EnvironmentFile=/home/ubuntu/brc-deploy/env/runtime-application.active.env
