# deploy/systemd/brc-runtime-signal-watcher.service.d/60-dry-run-audit-chain.conf
[Service]
# Refresh the non-executing runtime audit chain after each watcher tick.
# This uses local fixtures only; it does not call Tokyo APIs, exchange write,
# OrderLifecycle, withdrawals, transfers, secrets, live profile, or sizing.
ExecStartPost=/home/ubuntu/brc-deploy/venvs/brc-bnb-prelive-20260601/bin/python scripts/runtime_dry_run_audit_chain.py \
  --output-dir /home/ubuntu/brc-deploy/reports/runtime-signal-watcher/dry-run-audit-chain \
  --output-json /home/ubuntu/brc-deploy/reports/runtime-signal-watcher/runtime-dry-run-audit-chain.json
