# deploy/systemd/brc-runtime-signal-watcher.service.d/70-goal-status.conf
[Service]
# Refresh the main StrategyGroup runtime goal summary after each watcher tick.
# This is a read-only packet builder. It does not call FinalGate, Operation
# Layer, exchange write, OrderLifecycle, withdrawals, transfers, secrets, live
# profile, or sizing.
ExecStartPost=/home/ubuntu/brc-deploy/venvs/brc-bnb-prelive-20260601/bin/python scripts/build_strategygroup_runtime_goal_status.py \
  --report-dir /home/ubuntu/brc-deploy/reports/runtime-signal-watcher \
  --release-manifest /home/ubuntu/brc-deploy/app/current/.brc-release-manifest.json \
  --output-json /home/ubuntu/brc-deploy/reports/runtime-signal-watcher/strategygroup-runtime-goal-status.json
