# deploy/systemd/brc-runtime-signal-watcher.service.d/70-goal-status.conf
[Service]
# Retired as a final goal-status writer.
# The final StrategyGroup runtime goal summary is written by
# 80-product-state-refresh.conf after Candidate Pool, Daily Table, and Single
# Lane Packet are refreshed, using --candidate-pool-json as the control-plane
# authority. This drop-in intentionally has no ExecStartPost writer so legacy
# pilot/status scope alignment cannot overwrite the final Candidate Pool-backed
# goal status.
