# deploy/systemd/brc-runtime-signal-watcher.service.d/80-product-state-refresh.conf
[Service]
# Refresh Owner Console product-state packets after watcher and goal-status
# packets are current. This is a read-only API/readmodel refresh; it does not
# call FinalGate, Operation Layer, exchange write, OrderLifecycle, withdrawals,
# transfers, secrets, live profile, or sizing.
ExecStartPost=/home/ubuntu/brc-deploy/venvs/brc-bnb-prelive-20260601/bin/python scripts/refresh_strategygroup_runtime_product_state_packets.py \
  --api-base http://127.0.0.1:18080 \
  --output-dir /home/ubuntu/brc-deploy/reports/runtime-signal-watcher \
  --output-json /home/ubuntu/brc-deploy/reports/runtime-signal-watcher/product-state-refresh-packet.json \
  --collect-live-facts-before-refresh \
  --handoff-dir /home/ubuntu/brc-deploy/app/current/docs/current/strategy-group-handoffs \
  --env-file /home/ubuntu/brc-deploy/env/live-readonly.env \
  --live-facts-output /home/ubuntu/brc-deploy/reports/runtime-signal-watcher/strategy-group-live-facts-input.json \
  --refresh-chain-closure-status \
  --chain-closure-output-json /home/ubuntu/brc-deploy/reports/runtime-signal-watcher/runtime-execution-chain-closure-status.json \
  --refresh-live-closure-evidence \
  --live-closure-evidence-report-dir /home/ubuntu/brc-deploy/reports/runtime-signal-watcher \
  --live-closure-evidence-output-json /home/ubuntu/brc-deploy/reports/runtime-signal-watcher/runtime-live-closure-evidence.json \
  --live-closure-evidence-verification-output-json /home/ubuntu/brc-deploy/reports/runtime-signal-watcher/runtime-live-closure-evidence-verification.json \
  --live-closure-evidence-refresh-output-json /home/ubuntu/brc-deploy/reports/runtime-signal-watcher/runtime-live-closure-evidence-refresh.json \
  --label tokyo-runtime-signal-watcher
