# deploy/systemd/brc-runtime-signal-watcher.service.d/80-product-state-refresh.conf
[Service]
# Refresh PG current projections after watcher PG state is current.
# The Python sequence runner returns a step-level stdout summary but does not
# write a recurring sequence report file.
# This is a non-authority refresh. In watcher_tick_summary mode it publishes PG
# current projections only. Action-time modes may run ticket-bound preflight or
# handoff materializers, but this drop-in never calls FinalGate submit authority,
# Operation Layer submit, exchange write, OrderLifecycle, withdrawals, transfers,
# secrets, live profile, or sizing.
ExecStartPost=/home/ubuntu/brc-deploy/venvs/brc-bnb-prelive-20260601/bin/python scripts/run_server_product_state_refresh_sequence.py \
  --python /home/ubuntu/brc-deploy/venvs/brc-bnb-prelive-20260601/bin/python \
  --api-base http://127.0.0.1:18080 \
  --mode watcher_tick_summary
