# deploy/systemd/brc-runtime-signal-watcher.service.d/80-product-state-refresh.conf
[Service]
ExecStartPost=/usr/bin/time -f 'BRC_STAGE_MAX_RSS_KIB stage=projection_summary value=%M' /usr/bin/timeout --foreground --signal=TERM --kill-after=2s 8s /home/ubuntu/brc-deploy/app/current/.venv/bin/python scripts/run_server_product_state_refresh_sequence.py --python /home/ubuntu/brc-deploy/app/current/.venv/bin/python --api-base http://127.0.0.1:18080 --mode watcher_tick_summary
