# deploy/systemd/brc-runtime-signal-watcher.service.d/80-product-state-refresh.conf
[Service]
# Refresh Owner Console product-state packets after watcher and goal-status
# packets are current. This is a read-only API/readmodel refresh; it does not
# call FinalGate, Operation Layer, exchange write, OrderLifecycle, withdrawals,
# transfers, secrets, live profile, or sizing.
ExecStartPost=/home/ubuntu/brc-deploy/venvs/brc-bnb-prelive-20260601/bin/python scripts/refresh_strategygroup_runtime_product_state_packets.py \
  --api-base http://127.0.0.1:18080 \
  --output-dir /home/ubuntu/brc-deploy/reports/runtime-signal-watcher \
  --label tokyo-runtime-signal-watcher
