# deploy/systemd/brc-runtime-signal-watcher.service.d/80-product-state-refresh.conf
[Service]
# Refresh Owner Console product-state artifacts after watcher artifacts are current.
# The Python sequence runner records step-level results and still attempts the
# final Goal Status projection after optional readmodel failures.
# This is a non-authority refresh: it may create non-executing prepare evidence
# for one Candidate Pool action-time lane, but it does not call FinalGate,
# Operation Layer, exchange write, OrderLifecycle, withdrawals, transfers,
# secrets, live profile, or sizing.
ExecStartPost=/home/ubuntu/brc-deploy/venvs/brc-bnb-prelive-20260601/bin/python scripts/run_server_product_state_refresh_sequence.py \
  --python /home/ubuntu/brc-deploy/venvs/brc-bnb-prelive-20260601/bin/python \
  --api-base http://127.0.0.1:18080 \
  --report-dir /home/ubuntu/brc-deploy/reports/runtime-signal-watcher \
  --runtime-monitor-dir /home/ubuntu/brc-deploy/reports/runtime-monitor \
  --env-file /home/ubuntu/brc-deploy/env/live-readonly.env \
  --mode watcher_tick_summary \
  --output-json /home/ubuntu/brc-deploy/reports/runtime-signal-watcher/server-product-state-refresh-sequence.json
