# deploy/systemd/brc-runtime-signal-watcher.service.d/80-product-state-refresh.conf
[Service]
# Refresh Owner Console product-state artifacts after watcher artifacts are current.
# This also rebuilds goal-status after live closure evidence is refreshed, so a
# completed first bounded live order can be visible in the same watcher post-step.
# This is a read-only API/readmodel refresh; it does not call FinalGate,
# Operation Layer, exchange write, OrderLifecycle, withdrawals, transfers,
# secrets, live profile, or sizing.
ExecStartPost=/bin/sh -lc 'set -eu; cd /home/ubuntu/brc-deploy/app/current; PY=/home/ubuntu/brc-deploy/venvs/brc-bnb-prelive-20260601/bin/python; STATUS=/home/ubuntu/brc-deploy/reports/runtime-signal-watcher/latest-status.json; "$PY" scripts/build_mi_trial_admission_decision.py; "$PY" scripts/build_strategygroup_runtime_safety_state.py; "$PY" scripts/validate_runtime_candidate_universe_coverage.py "$STATUS"; "$PY" scripts/build_strategy_live_candidate_pool.py --runtime-active-monitor-json "$STATUS"; "$PY" scripts/validate_strategy_live_candidate_pool.py output/runtime-monitor/latest-strategy-live-candidate-pool.json; "$PY" scripts/build_daily_live_enablement_table.py --candidate-pool-json output/runtime-monitor/latest-strategy-live-candidate-pool.json; "$PY" scripts/validate_daily_live_enablement_table.py output/runtime-monitor/latest-daily-live-enablement-table.json; "$PY" scripts/build_single_lane_task_packet.py; "$PY" scripts/validate_single_lane_task_packet.py output/runtime-monitor/latest-single-lane-task-packet.json'
ExecStartPost=/home/ubuntu/brc-deploy/venvs/brc-bnb-prelive-20260601/bin/python scripts/refresh_strategygroup_runtime_product_state_artifacts.py \
  --api-base http://127.0.0.1:18080 \
  --output-dir /home/ubuntu/brc-deploy/reports/runtime-signal-watcher \
  --output-json /home/ubuntu/brc-deploy/reports/runtime-signal-watcher/product-state-refresh-packet.json \
  --collect-live-facts-before-refresh \
  --handoff-dir /home/ubuntu/brc-deploy/app/current/docs/current/strategy-group-handoffs \
  --env-file /home/ubuntu/brc-deploy/env/live-readonly.env \
  --live-facts-output /home/ubuntu/brc-deploy/reports/runtime-signal-watcher/strategy-group-live-facts-input.json \
  --refresh-chain-closure-status \
  --chain-closure-output-json /home/ubuntu/brc-deploy/reports/runtime-signal-watcher/runtime-execution-chain-closure-status.json \
  --refresh-live-closure-evidence \
  --live-closure-evidence-report-dir /home/ubuntu/brc-deploy/reports/runtime-signal-watcher \
  --live-closure-evidence-output-json /home/ubuntu/brc-deploy/reports/runtime-signal-watcher/runtime-live-closure-evidence.json \
  --live-closure-evidence-verification-output-json /home/ubuntu/brc-deploy/reports/runtime-signal-watcher/runtime-live-closure-evidence-verification.json \
  --live-closure-evidence-refresh-output-json /home/ubuntu/brc-deploy/reports/runtime-signal-watcher/runtime-live-closure-evidence-refresh.json \
  --refresh-goal-status \
  --goal-status-output-json /home/ubuntu/brc-deploy/reports/runtime-signal-watcher/strategygroup-runtime-goal-status.json \
  --release-manifest /home/ubuntu/brc-deploy/app/current/.brc-release-manifest.json \
  --label tokyo-runtime-signal-watcher
