# deploy/systemd/brc-runtime-signal-watcher.service.d/85-action-time-refresh-if-needed.conf
[Service]
ExecStartPost=/usr/bin/time -f 'BRC_STAGE_MAX_RSS_KIB stage=action_time_refresh value=%M' /usr/bin/timeout --foreground --signal=TERM --kill-after=3s 37s /home/ubuntu/brc-deploy/app/current/.venv/bin/python scripts/run_server_product_state_refresh_sequence.py --python /home/ubuntu/brc-deploy/app/current/.venv/bin/python --env-file /home/ubuntu/brc-deploy/env/live-readonly.env --mode action_time_if_needed
