# deploy/systemd/brc-runtime-signal-watcher.service.d/90-resume-dispatcher-after-refresh.conf
[Service]
ExecStartPost=/usr/bin/time -f 'BRC_STAGE_MAX_RSS_KIB stage=resume_dispatcher value=%M' /usr/bin/timeout --foreground --signal=TERM --kill-after=3s 42s /home/ubuntu/brc-deploy/app/current/.venv/bin/python scripts/runtime_signal_watcher_resume_dispatcher.py --identity-source pg_ticket --api-base http://127.0.0.1:18080 --label tokyo-runtime-signal-watcher --execute-preflight --preflight-timeout-seconds 30 --execute-operation-layer-submit --operation-layer-submit-mode from_submit_mode_decision --production-submit-execution-policy armed
