# deploy/systemd/brc-runtime-signal-watcher.service.d/90-resume-dispatcher-after-refresh.conf
[Service]
# Dispatcher consumes PG Action-Time Ticket identity only.
# TEMPORARY(L2-L9-closure): this branch arms a narrow ENTRY + SL + TP1 live
# aperture to validate real exchange behavior. Delete this submit mode and
# return this drop-in to preflight-only after the full post-submit protection,
# reconciliation, settlement, and review chain is closed.
ExecStartPost=/home/ubuntu/brc-deploy/venvs/brc-bnb-prelive-20260601/bin/python scripts/runtime_signal_watcher_resume_dispatcher.py \
  --identity-source pg_ticket \
  --api-base http://127.0.0.1:18080 \
  --label tokyo-runtime-signal-watcher \
  --execute-preflight \
  --execute-operation-layer-submit \
  --operation-layer-submit-mode temp_tiny_live_protected_submit
