# deploy/systemd/brc-runtime-signal-watcher.service.d/90-resume-dispatcher-after-refresh.conf
[Service]
ExecStartPost=/usr/bin/time -f 'BRC_STAGE_MAX_RSS_KIB stage=resume_dispatcher value=%M' /usr/bin/timeout --foreground --signal=TERM --kill-after=3s 42s /home/ubuntu/brc-deploy/app/current/.venv/bin/python scripts/run_action_time_dispatch_command_once.py --require-database-url --worker-id tokyo-runtime-signal-watcher --production-submit-execution-policy armed
