# deploy/systemd/brc-runtime-signal-watcher.service.d/90-resume-dispatcher-after-refresh.conf
[Service]
# Dispatcher must consume the resume/readiness pack rebuilt by
# 80-product-state-refresh.conf after Candidate Pool lane materialization.
# This keeps same-tick fresh lanes from being delayed by one watcher cycle.
ExecStartPost=/home/ubuntu/brc-deploy/venvs/brc-bnb-prelive-20260601/bin/python scripts/runtime_signal_watcher_resume_dispatcher.py \
  --report-dir /home/ubuntu/brc-deploy/reports/runtime-signal-watcher \
  --resume-dispatch-json /home/ubuntu/brc-deploy/reports/runtime-signal-watcher/resume-dispatch-artifact.json \
  --api-base http://127.0.0.1:18080 \
  --owner-operator-id owner \
  --owner-confirmation-reference owner-authorized-runtime-resume-dispatch \
  --execute-preflight \
  --execute-operation-layer-submit \
  --operation-layer-submit-mode real_gateway_action \
  --execute-post-submit-finalize
