# deploy/systemd/brc-runtime-monitor.timer
[Unit]
Description=Run BRC server-side readonly runtime monitor every 10 minutes

[Timer]
OnBootSec=120s
OnUnitActiveSec=10min
AccuracySec=30s
Persistent=true
Unit=brc-runtime-monitor.service

[Install]
WantedBy=timers.target
