# deploy/systemd/brc-trading-kernel.slice
[Unit]
Description=BRC trading kernel resource boundary

[Slice]
CPUAccounting=true
MemoryAccounting=true
CPUQuota=100%
MemoryMax=1G
TasksMax=128
