# deploy/systemd/brc-runtime-db-retention.timer
[Unit]
Description=Run BRC runtime-control-state PG retention cleanup daily

[Timer]
OnCalendar=*-*-* 03:20:00 Asia/Shanghai
RandomizedDelaySec=10min
Persistent=true
Unit=brc-runtime-db-retention.service

[Install]
WantedBy=timers.target
