# deploy/systemd/brc-trading-kernel-lifecycle-worker.timer
[Unit]
Description=Run the BRC trading kernel lifecycle worker every two seconds

[Timer]
OnBootSec=16s
OnUnitActiveSec=2s
AccuracySec=100ms
Persistent=false
Unit=brc-trading-kernel-lifecycle-worker.service

[Install]
WantedBy=timers.target
