# deploy/systemd/brc-runtime-signal-watcher.service
[Unit]
Description=BRC runtime signal watcher tick
Wants=network-online.target
After=network-online.target brc-owner-console-backend.service

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/brc-deploy/app/current
EnvironmentFile=/home/ubuntu/brc-deploy/env/live-readonly.env
# Optional bounded-experiment overlay. When present, this can raise the runtime
# ceiling from read-only to order-capable without changing credentials or scope.
EnvironmentFile=-/home/ubuntu/brc-deploy/env/runtime-order-capable.env
Environment=BRC_STRATEGYGROUP_MAX_SYMBOLS=3
Environment=BRC_STRATEGYGROUP_STALE_AFTER_SECONDS=180
EnvironmentFile=-/home/ubuntu/brc-deploy/env/runtime-signal-watcher.env
ExecStartPre=/home/ubuntu/brc-deploy/venvs/brc-bnb-prelive-20260601/bin/python scripts/fetch_binance_usdm_public_facts.py \
  --require-database-url
ExecStart=/home/ubuntu/brc-deploy/venvs/brc-bnb-prelive-20260601/bin/python scripts/runtime_signal_watcher_tick.py \
  --env-file /home/ubuntu/brc-deploy/env/live-readonly.env \
  --api-base http://127.0.0.1:18080 \
  --source live_market \
  --strategy-source live_market \
  --max-iterations 1 \
  --cycle-timeout-seconds 180 \
  --status-stale-after-seconds 900 \
  --require-database-url \
  --allow-prepare-records \
  --allow-arm-preview \
  --allow-attempt-policy-prepare \
  --allow-disabled-smoke \
  --allow-standing-operation-layer-evidence-prep \
  --strategy-family-id CPM-RO-001 \
  --strategy-family-id MPG-001 \
  --strategy-family-id MI-001 \
  --strategy-family-id SOR-001 \
  --strategy-family-id BRF2-001 \
  --label tokyo-runtime-signal-watcher
NoNewPrivileges=true
PrivateTmp=true
Nice=10
CPUAccounting=true
CPUQuota=60%
IOAccounting=true
IOSchedulingClass=best-effort
IOSchedulingPriority=6
