# deploy/systemd/brc-trading-kernel-lifecycle-worker.service
[Unit]
Description=BRC trading kernel bounded lifecycle worker tick
After=network-online.target
Wants=network-online.target
ConditionPathExists=/etc/brc/trading-kernel.env

[Service]
Type=oneshot
User=brc
Group=brc
WorkingDirectory=/opt/brc/current
EnvironmentFile=/etc/brc/trading-kernel.env
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/opt/brc/current/.venv/bin/python scripts/trading_kernel/run_command_worker_once.py \
    --venue-factory ${TRADING_KERNEL_VENUE_FACTORY} \
    --worker-role lifecycle \
    --worker-id ${TRADING_KERNEL_LIFECYCLE_WORKER_ID} \
    --lease-ms ${TRADING_KERNEL_LEASE_MS} \
    --timeout-seconds ${TRADING_KERNEL_TIMEOUT_SECONDS} \
    --idle-poll-interval-ms ${TRADING_KERNEL_IDLE_POLL_INTERVAL_MS}
TimeoutStartSec=30
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict

[Install]
WantedBy=multi-user.target
