# deploy/systemd/brc-runtime-signal-watcher.timer
[Unit]
Description=Run BRC runtime signal watcher on bounded cadence

[Timer]
OnActiveSec=90s
OnUnitActiveSec=3min
AccuracySec=15s
Persistent=false
Unit=brc-runtime-signal-watcher.service

[Install]
WantedBy=timers.target
