# deploy/systemd/brc-trading-kernel-reconciliation-worker.timer
[Unit]
Description=Run the BRC trading kernel reconciliation worker every two seconds

[Timer]
OnBootSec=17s
OnUnitActiveSec=2s
AccuracySec=100ms
Persistent=false
Unit=brc-trading-kernel-reconciliation-worker.service

[Install]
WantedBy=timers.target
