# deploy/systemd/brc-runtime-signal-watcher.timer
[Unit]
Description=Run BRC runtime signal watcher on bounded cadence

[Timer]
OnActiveSec=5min
OnUnitActiveSec=10min
AccuracySec=30s
Persistent=false
Unit=brc-runtime-signal-watcher.service

[Install]
WantedBy=timers.target
