# deploy/systemd/brc-trading-kernel-observation-worker.timer
[Unit]
Description=Check for one due BRC trading kernel observation every five seconds

[Timer]
OnBootSec=14s
OnUnitActiveSec=5s
AccuracySec=100ms
Persistent=false
Unit=brc-trading-kernel-observation-worker.service

[Install]
WantedBy=timers.target
