# deploy/systemd/brc-ticket-lifecycle-maintenance.service.d/50-runtime-release-identity.conf
[Service]
Environment=BRC_RUNTIME_HEAD=__BRC_CANDIDATE_SHA__
