# deploy/systemd/brc-owner-console-backend.service.d/10-runtime-bound.conf
[Service]
ExecStart=
ExecStart=/home/ubuntu/brc-deploy/app/current/.venv/bin/python -m src.main
