# deploy/systemd/brc-owner-console-backend.service.d/30-runtime-order-capable-identity.conf
[Service]
# Non-secret identity values remain server-owned in this environment file.
# The repository owns only the process-boundary wiring.
EnvironmentFile=-/home/ubuntu/brc-deploy/env/runtime-order-capable.env
