# deploy/systemd/brc-owner-console-backend.service.d/40-runtime-stability.conf
[Unit]
After=docker.service network-online.target
Wants=docker.service network-online.target
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
ExecStartPre=/usr/bin/timeout --foreground --signal=TERM --kill-after=2s 8s /home/ubuntu/brc-deploy/app/current/.venv/bin/python /home/ubuntu/brc-deploy/control-plane/check_runtime_postgres_ready.py --require-database-url --timeout-seconds 8
Restart=on-failure
RestartSec=15s
TimeoutStopSec=20s
KillMode=control-group
