# deploy/systemd/brc-owner-console-backend.service.d/45-persistent-config-path.conf
[Service]
# Configuration storage is machine state, not mutable release content.
Environment=CONFIG_DB_PATH=/var/lib/brc-runtime/config/v3_dev.db
