# deploy/systemd/brc-owner-console-backend.service.d/60-runtime-db-application.conf
[Service]
# Created only inside the Writer Fence deployment transaction.  It overlays
# the legacy server-owned exchange/notification environment with the least-
# privilege PostgreSQL application identity.
EnvironmentFile=/home/ubuntu/brc-deploy/env/runtime-application.active.env
